# homepage.container
[Container]
Image=ghcr.io/gethomepage/homepage:latest
# 网络模式 - 使用 host 网络模式
Network=host
# 环境变量
Environment=PUID=1000
Environment=PGID=1000
Environment=TZ=Asia/Shanghai
Environment=HOMEPAGE_ALLOWED_HOSTS=*
# 卷挂载
Volume=/data/appdata/homepage:/app/config:Z
Volume=/data/appdata/homepage/icons:/app/public/icons:Z
Volume=/run/user/1000/podman/podman.sock:/var/run/docker.sock
# 如果需要切换回桥接模式，可以注释上面的 Network=host 并取消注释下面的端口映射
# PublishPort=3000:3000

[Service]
# 重启策略
Restart=unless-stopped
# 可选：添加 systemd 依赖控制
Wants=network-online.target
After=network-online.target
[Install]
# 开机自启动
WantedBy=default.target
