# power.container
[Container]
# 使用自定义镜像或直接运行 Python 应用
Image=docker.io/library/python:3
# 工作目录
WorkingDir=/usr/src/myapp
# 环境变量
Environment=PUID=1000
Environment=PGID=1000
Environment=TZ=Asia/Shanghai
# 端口映射
PublishPort=5008:5008
# 卷挂载 - 修正 SELinux 标签
Volume=/data/www:/usr/src/myapp:z
# 启动命令 - 使用启动脚本
Exec=/usr/src/myapp/start.sh

[Service]
# 重启策略
Restart=unless-stopped
# 可选：添加 systemd 依赖控制
Wants=network-online.target
After=network-online.target
[Install]
# 开机自启动
WantedBy=default.target
